# supervisord.conf
[supervisord]
nodaemon=true

[program:speech2txtbot]
command=python speech2textbot.py
